# container/dockerfiles/jetty9/supervisord-jetty9.conf
[program:jetty9]
command=/start-jetty9.sh
numprocs=1
autostart=true
user=jetty
autorestart=true
